# gowwwredirect.conf
description "www redirection server"
 
start on filesystem or runlevel [2345]
stop on run level [!2345]
 
env USER='root'
env APP_DIR='/usr/local/sbin/'
env APP_EXEC='gowwwredirect'
 
exec start-stop-daemon --start --chuid ${USER} --chdir ${APP_DIR} --exec ${APP_DIR}${APP_EXEC}